# networking-bagpipe.service
[Unit]
Description=BaGPipe-BGP server
After=syslog.target network.target

[Service]
Type=simple
User=neutron
ExecStart=/usr/bin/bagpipe-bgp
ExecStop=/usr/bin/bagpipe-bgp-cleanup
KillMode=process
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
